# gitea/rootfs/etc/services.d/gitea/run
#!/bin/sh
set -e

echo "==> Gitea addon starting..."
echo "==> ENV vars:"
env | grep GITEA || true
echo "==> End ENV vars"
